# deploy/backups-rotate.service
[Unit]
Description=invoke backups-rotate

# keep backups during two months
[Service]
ExecStart=/root/venv-for-backups-rotate/bin/backups-rotate
