# deploy/backups-rotate.service
[Unit]
Description=invoke backups-rotate

# keep backups during two months
[Service]
ExecStart=/bin/bash -c backups-rotate
