# deploy/r2lab-backup.timer
[Unit]
Description=Hourly R2Lab database backup

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
